# ouvretaferme/main/conf/crontab/prod
############################
# www.dev-ouvretaferme.org #
############################

* * * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme website/cron/certbot:ping
* * * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme website/cron/certbot:rewrite
*/5 * * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme session/cron/sessionClean
0 * * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme mail/cron/mail
0 * * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme selling/cron/invoice
0 * * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme storage/cron/buffer
0 0 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme user/cron/seniority
0 12 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme mail/cron/clean
0 2 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme media/cron/clean
0 3 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme hr/cron/calculate
0 4 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme farm/cron/demo
0 6 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme cron/sitemap
0 6 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme user/cron/close
0 6 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme user/cron/log
17 7 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme dev/cron/minifyClean
19 5 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme selling/cron/clean
28 4 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme website/cron/certbot:clean
3 */12 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme website/cron/certbot:certificate
47 5 * * * php /var/www/otf/framework/lime.php -e prod -a ouvretaferme dev/cron/error